# deploy/shining-planet.service
[Unit]
Description=Shining Planet API
After=network.target

[Service]
Type=simple
User=shareholder
Group=shareholder
WorkingDirectory=/home/shareholder/shining-planet/backend
Environment=PATH=/home/shareholder/shining-planet/venv/bin
ExecStartPre=/home/shareholder/shining-planet/venv/bin/python3 -c "from app.database import init_db; init_db()"
# 1C2G 机器：单 worker，降低内存占用；AI 请求单独拉长超时
ExecStart=/home/shareholder/shining-planet/venv/bin/gunicorn app.main:app -w 1 -k uvicorn.workers.UvicornWorker -b 127.0.0.1:8020 --timeout 180 --max-requests 200 --max-requests-jitter 40
Restart=always
RestartSec=3
MemoryMax=450M
MemoryHigh=380M

[Install]
WantedBy=multi-user.target
